# yamlspecs/nvidia-installer.service.in
[Unit]
Description=run the nvidia installer if not already installed. It will also run nvidia-smi to build /dev/nvidia
After=network.target rsyslog.service

[Service]
Type=oneshot
ExecStart=%PATH%/installer.sh
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
